# config/crontab
* * * * * root echo "Hello, world!" >> /tmp/cron_output.txt
